# deploy/aliyun-lightserver/relay.service
[Unit]
Description=Talk with Buffett ASR Relay
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/talk-with-buffett
Environment=NODE_ENV=production
Environment=PORT=3001
EnvironmentFile=/root/talk-with-buffett/.env.relay
ExecStart=/usr/bin/npm run relay:start
Restart=always
RestartSec=5
User=root
Group=root

[Install]
WantedBy=multi-user.target
